# lib/Docker/php8.1/supervisor/supervisor.program.conf
[program:messenger]
command=php console pipeline:run-single
user=root
startsecs=2
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d
directory=/var/app/bin
